# mldev.setup/files/systemd/system/jupyter.service
[Unit]
Description=Jupyter notebook
After=network.target

[Service]
WorkingDirectory=/root/work/notebooks
ExecStart=/usr/local/bin/jupyter notebook
